# packages/pglift/pglift-0.7.0.tar.gz/pglift-0.7.0/src/pglift/data/systemd/postgresql-backup.timer
[Unit]
Description=Backup %i PostgreSQL database instance
Requires=postgresql@%i.service
After=postgresql@%i.service

[Timer]
OnCalendar={calendar}
Persistent=true

[Install]
WantedBy=timers.target
