# packages/pglift/pglift-0.7.0-py3-none-any.whl/pglift/data/systemd/postgresql-backup.service
[Unit]
Description=Backup %i PostgreSQL database instance
After=postgresql@%i.service

[Service]
Type=oneshot
{environment}
ExecStart={python} -m pglift.backup %i
